# build/system/etc/systemd/system/bthelper@.service.d/01-add-delay.conf
[Service]
ExecStartPre=
ExecStartPre=/bin/sleep 2