# build/system/etc/systemd/system/bthelper@.service.d/01-add-delay.conf
[Unit]
After=hciuart.service bluetooth.service
Before=

[Service]
ExecStartPre=
ExecStartPre=/bin/sleep 5
