# DGTCentaurMods/etc/centaurmods-web.service
[Unit]
Description=Centaur Web interface
After=network.target

[Service]
ExecStart=flask run
WorkingDirectory=/home/pi/DGTCentaur/DGTCentaurMods/web/
Environment=
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
