# DGTCentaurMods/etc/systemd/system/bluetooth.service.d/01-disable-sap-plugin.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --compat --noplugin=sap,battery