# DGTCentaurMods/etc/systemd/system/bluetooth.service.d/02-sdptool-add-SP.conf
[Service]
ExecStartPost=
ExecStartPost=/usr/bin/sdptool add SP