# build/system/etc/systemd/system/centaur.service
[Unit]
Description=Centaur board original software
After=network.target
# Start other related services
#Requires=centaur-helper.service

[Service]
ExecStart=/home/pi/centaur/centaur
WorkingDirectory=/home/pi/centaur/
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target
